# supervisord.conf
[supervisord]

[program:postgres]
user=postgres
command=/usr/lib/postgresql/9.3/bin/postgres -D /var/lib/postgresql/9.3/main -c config_file=/etc/postgresql/9.3/main/postgresql.conf
redirect_stderr=true
autorestart=true

[program:apache2]
command=/bin/bash -c "source /etc/apache2/envvars && exec /usr/sbin/apache2 -DFOREGROUND"


[program:galaxy]
command=/bin/sh run.sh
redirect_stderr=true
autorestart=true
